# os-config/quadlets/autorestic-restore.container
[Unit]
Description=Autorestic Restore
After=network-online.target fetch-bootstrap-secrets.service
Wants=network-online.target
Requires=fetch-bootstrap-secrets.service
Before=freshrss.service planka.service nitter.service redlib.service

[Container]
Image=docker.io/cupcakearmy/autorestic:1.8.3@sha256:0a826ea59ca452dadae50e618c0a14800450c920d22da0e35cb23e7a8be36d1f
Exec=sh -c '\
  for loc in freshrss planka; do \
    dir="/opt/data/$loc"; \
    [ ! -d "$dir" ] || [ -z "$(ls -A "$dir" 2>/dev/null)" ] && autorestic restore -l "$loc" --to / --force; \
  done; \
  chown -R 1000:1000 /opt/data/planka/favicons /opt/data/planka/avatars /opt/data/planka/backgrounds /opt/data/planka/attachments 2>/dev/null; \
  true'
Network=host
Volume=/etc/rclone/rclone.conf:/root/.config/rclone/rclone.conf:ro,Z
Volume=/root/.config/autorestic/autorestic.yml:/root/.autorestic.yml:ro,Z
Volume=/opt/data:/opt/data:z
Volume=/run/secrets/restic-password:/run/secrets/restic-password:ro,z
Environment=AUTORESTIC_CONFIG=/root/.autorestic.yml
Environment=RESTIC_PASSWORD_FILE=/run/secrets/restic-password
EnvironmentFile=/etc/bucket/config
DropCapability=ALL
AddCapability=CHOWN
AddCapability=DAC_OVERRIDE
AddCapability=DAC_READ_SEARCH
AddCapability=FOWNER
PodmanArgs=--security-opt=no-new-privileges

[Service]
Type=oneshot
ExecCondition=/bin/bash -c '\
  for dir in /opt/data/freshrss /opt/data/planka; do \
    [ ! -d "$dir" ] || [ -z "$(ls -A "$dir" 2>/dev/null)" ] && exit 0; \
  done; \
  exit 1'
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
