# os-config/quadlets/traefik.container
[Unit]
Description=Traefik Reverse Proxy
After=podman.socket fetch-config.service network-online.target
Wants=network-online.target
Requires=fetch-config.service podman.socket

[Container]
Image=docker.io/traefik:v3.6.7
PublishPort=80:80
Volume=/run/podman/podman.sock:/var/run/docker.sock:z
Volume=/root/.config/traefik:/etc/traefik:ro,Z
SecurityLabelType=container_runtime_t

[Service]
Restart=always

[Install]
WantedBy=multi-user.target
