# os-config/systemd/app-restart@.service
[Unit]
Description=Restart %i on secret change

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl restart %i.service
