# os-config/systemd/freshrss.service
[Unit]
Description=FreshRSS
After=traefik.service infisical-agent.service autorestic-restore.service ghcr-login.service freshrss-feeds-network.service
Wants=autorestic-restore.service
Requires=ghcr-login.service freshrss-feeds-network.service
ConditionPathExists=/opt/secrets/freshrss.env

[Service]
Type=oneshot
WorkingDirectory=/opt/containers/freshrss
ExecStart=/usr/bin/podman compose --env-file /opt/env/freshrss.env up -d
ExecStartPost=/opt/scripts/freshrss-cache-warm.sh
ExecStop=/usr/bin/podman compose down
RemainAfterExit=yes
TimeoutStartSec=900

[Install]
WantedBy=multi-user.target
