# os-config/systemd/ghcr-token.path
[Unit]
Description=Watch for GHCR token file

[Path]
PathExists=/opt/secrets/ghcr-token

[Install]
WantedBy=multi-user.target
