# os-config/systemd/languagetool.path
[Unit]
Description=Watch LanguageTool for changes

[Path]
PathChanged=/opt/containers/languagetool/
Unit=app-restart@languagetool.service

[Install]
WantedBy=multi-user.target
