# os-config/systemd/nitter.path
[Unit]
Description=Watch Nitter for changes

[Path]
PathChanged=/opt/containers/nitter/
PathModified=/opt/secrets/nitter.env
PathModified=/opt/secrets/nitter-sessions.jsonl
Unit=app-restart@nitter.service

[Install]
WantedBy=multi-user.target
