# os-config/systemd/redlib.path
[Unit]
Description=Watch Redlib for changes

[Path]
PathChanged=/opt/containers/redlib/
PathModified=/opt/secrets/redlib.env
Unit=app-restart@redlib.service

[Install]
WantedBy=multi-user.target
